# systemd/linux-management-prod.service
[Unit]
Description=Linux Management System - Production Environment
Documentation=https://github.com/Kensan196948G/Linux-Management-System
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=kensan
Group=kensan
WorkingDirectory=/mnt/LinuxHDD/Linux-Management-Systm

# 環境変数
Environment="ENV=prod"
Environment="PYTHONUNBUFFERED=1"
Environment="PATH=/home/kensan/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

# .env ファイルから静的環境変数を読み込む
EnvironmentFile=/mnt/LinuxHDD/Linux-Management-Systm/.env

# 起動前: IPアドレスを自動検出して .env.runtime に書き出す
ExecStartPre=/mnt/LinuxHDD/Linux-Management-Systm/scripts/detect-ip.sh

# .env.runtime（detect-ip.sh が生成）を読み込む
EnvironmentFile=-/mnt/LinuxHDD/Linux-Management-Systm/.env.runtime

# 実行コマンド（本番: gunicorn + uvicorn worker）
ExecStart=/home/kensan/.local/bin/gunicorn \
    backend.api.main:app \
    --bind 0.0.0.0:8000 \
    --workers 2 \
    --worker-class uvicorn.workers.UvicornWorker \
    --timeout 120 \
    --keep-alive 5 \
    --max-requests 1000 \
    --max-requests-jitter 50 \
    --access-logfile - \
    --error-logfile - \
    --log-level warning

# 再起動ポリシー（本番環境）
Restart=always
RestartSec=10s

# タイムアウト
TimeoutStartSec=60
TimeoutStopSec=30

# リソース制限
LimitNOFILE=65536

# セキュリティ強化
PrivateTmp=yes
NoNewPrivileges=yes

# ログ設定
StandardOutput=journal
StandardError=journal
SyslogIdentifier=linux-management-prod

[Install]
WantedBy=multi-user.target
